# backend/scripts/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=mongod --bind_ip_all

[program:cron]
command=cron -f
